# telepathy-mission-control/server/mission-control-5.service.in
[Unit]
Description=Telepathy mission control daemon
After=dbus.socket pre-user-session.target
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/invoker --type=generic @libexecdir@/mission-control-5
Type=dbus
BusName=org.freedesktop.Telepathy.MissionControl5

[Install]
WantedBy=user-session.target
